# crontab
# Run the backup script every minute
* * * * * /app/backup.sh >> /var/log/cron.log 2>&1
* * * * * /usr/bin/python3 /app/health_monitor.py >> /var/log/health_monitor.log 2>&1
# New line character required!
